# deployment/nrd-backend.service
[Unit]
Description=NRD Dashboard Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nrd-dashboard
Environment="PATH=/home/ubuntu/nrd-dashboard/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=-/home/ubuntu/nrd-dashboard/backend/.env
ExecStart=/home/ubuntu/nrd-dashboard/.venv/bin/python -m uvicorn backend.server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
